# imageroot/crowdsec-upgrade-hub.timer
[Unit]
Description=Timer of crowdsec hub update

[Timer]
OnActiveSec=4 hours
OnUnitInactiveSec=15 days

[Install]
WantedBy=timers.target
